# ansible/files/alok.service
[Install]
WantedBy=multi-user.target

[Unit]
After=multi-user.target

[Service]
Type=simple
WorkingDirectory=/usr/share/alok
User=rgl
Group=rgl
ExecStart=/bin/java \
    -Xms128M -Xmx128M \
    -jar /usr/share/alok/alok.jar \
    /home/rgl/alok.txt
KillMode=mixed
KillSignal=SIGTERM
SuccessExitStatus=143

TimeoutSec=1m
RestartSec=60
Restart=always
